# crontab
#run python script every minute
* * * * * python3 /app/scraper.py > /proc/1/fd/1 2>/proc/1/fd/2
